# systemd/volumectl.service
[Unit]
Description=Volume daemon
After=pulseaudio.service

[Service]
ExecStart=/usr/bin/volumectl -D
Restart=always

[Install]
WantedBy=basic.target
